# images/Custom/tomcat_ssh/supervisor_tomcat.conf
[program:tomcat]
directory=/
command=/tomcat/apache-tomcat-8.0.23/bin/catalina.sh run
user=root
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log
